# Library/LaunchAgents/nl.carosi.sleepwatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>de.bernhard-baehr.sleepwatcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/Cellar/sleepwatcher/2.2.1/sbin/sleepwatcher</string>
		<string>-V</string>
		<string>-s /Users/carosi/.sleepwatcher/sleep.sh</string>
		<string>-w /Users/carosi/.sleepwatcher/wakeup.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/carosi/Library/Logs/Sleepwatcher.log</string>
	<key>StandardOutPath</key>
	<string>/Users/carosi/Library/Logs/Sleepwatcher.log</string>
</dict>
</plist>
